# src/composite-metrics@.path.in
[Unit]
Description=composite_metrics activation for %I
Requires=malamute.service
After=malamute.service
PartOf=bios.target

[Path]
Environment="prefix=@prefix@"
PathChanged=@sysconfdir@/@PACKAGE@/%i.cfg
Unit=composite-metrics@%i.service

[Install]
WantedBy=bios.target
