# src/composite-metrics@.service.in
[Unit]
Description=composite_metrics service for %I
Requires=malamute.service
After=malamute.service
Conflicts=composite-metrics@%i.path
PartOf=bios.target

[Service]
Type=simple
User=nobody
EnvironmentFile=-/usr/share/bios/etc/default/bios
EnvironmentFile=-/usr/share/bios/etc/default/bios__%n.conf
EnvironmentFile=-/etc/default/bios
EnvironmentFile=-/etc/default/bios__composite-metrics__%i.conf
ExecStart=@prefix@/bin/composite-metrics /etc/@PACKAGE@/%i.cfg

[Install]
WantedBy=bios.target
