# src/composite-metrics@.path.in
[Unit]
Description=composite_metrics activation for %I
After=malamute.service

[Path]
Environment="prefix=@prefix@"
PathChanged=@sysconfdir@/@PACKAGE@/%i.cfg
Unit=composite-metrics@%i.service

[Install]
WantedBy=multi-user.target
