# src/composite-metrics-configurator.service.in
[Unit]
Description=composite-metrics-configurator service
Requires=malamute.service
After=malamute.service
PartOf=bios.target

[Service]
Type=simple
User=bios
EnvironmentFile=-/usr/share/bios/etc/default/bios
EnvironmentFile=-/usr/share/bios/etc/default/bios__%n.conf
EnvironmentFile=-/etc/default/bios
EnvironmentFile=-/etc/default/bios__composite-metrics__%i.conf
Environment="prefix=@prefix@"
ExecStart=@prefix@/bin/composite-metrics-configurator
Restart=always

[Install]
WantedBy=bios.target
