# src/composite_metrics@.service.in
[Unit]
Description=composite_metrics service for %I
After=network.target

[Service]
Type=simple
ExecStart=@prefix@/bin/composite_metrics @sysconfdir@/@PACKAGE@/%i.cfg

[Install]
WantedBy=multi-user.target
