# src/composite_metrics@.service.in
[Unit]
Description=composite_metrics service for %I
After=malamute.service

[Service]
Type=simple
Environment="prefix=@prefix@"
ExecStart=@prefix@/bin/composite_metrics @sysconfdir@/@PACKAGE@/%i.cfg

[Install]
WantedBy=multi-user.target
