# src/composite_metrics@.service.in
[Unit]
Description=composite_metrics service for %I
After=malamute.service

[Service]
Type=simple
ExecStart=@prefix@/bin/composite_metrics @prefix@/etc/@PACKAGE@/%i.cfg

[Install]
WantedBy=multi-user.target
