# src/dc_th.timer
[Unit]
Description=DC TH averaging configuration trigger

[Timer]
OnBootSec=5min
OnUnitActiveSec=15min
