# src/dc_th.service.in
[Unit]
Description=Enables computation of Temperature & Humidity in DC
After=mysql.service bios-db-init

[Service]
Type=oneshot
ExecStart=@prefix@/bin/dc_th_enable


[Install]
WantedBy=multi-user.target
