# src/fty-metric-composite-configurator.service.in
[Unit]
Description=fty-metric-composite-configurator service
Requires=malamute.service
After=malamute.service
PartOf=bios.target

[Service]
Type=simple
User=bios
EnvironmentFile=-/usr/share/bios/etc/default/bios
EnvironmentFile=-/usr/share/bios/etc/default/fty__%n.conf
EnvironmentFile=-/etc/default/bios
EnvironmentFile=-/etc/default/fty__fty-metric-composite__%i.conf
Environment="prefix=@prefix@"
ExecStart=@prefix@/bin/fty-metric-composite-configurator 
Restart=always

[Install]
WantedBy=bios.target
