# src/fty-metric-composite@.service.in
[Unit]
Description=fty-metric-composite service for %I
Requires=malamute.service
After=malamute.service
Conflicts=fty-metric-composite@%i.path
PartOf=bios.target

[Service]
Type=simple
User=bios
Restart=always
EnvironmentFile=-@prefix@/share/bios/etc/default/bios
EnvironmentFile=-@prefix@/share/bios/etc/default/bios__%n.conf
EnvironmentFile=-@prefix@/share/fty/etc/default/fty
EnvironmentFile=-@prefix@/share/fty/etc/default/fty__%n.conf
EnvironmentFile=-@sysconfdir@/default/bios
EnvironmentFile=-@sysconfdir@/default/bios__%n.conf
EnvironmentFile=-@sysconfdir@/default/fty
EnvironmentFile=-@sysconfdir@/default/fty__%n.conf
EnvironmentFile=-@sysconfdir@/default/fty__fty-metric-composite__%i.conf
Environment="prefix=@prefix@"
ExecStart=@prefix@/bin/fty-metric-composite /var/lib/fty/fty-metric-composite/%i.cfg
Restart=always

[Install]
WantedBy=bios.target
