# src/fty-metric-composite@.service.in
[Unit]
Description=fty-metric-composite service for %I
Requires=malamute.service
After=malamute.service
Conflicts=fty-metric-composite@%i.path
PartOf=bios.target

[Service]
Type=simple
User=bios
EnvironmentFile=-/usr/share/bios/etc/default/bios
EnvironmentFile=-/usr/share/bios/etc/default/fty__%n.conf
EnvironmentFile=-/etc/default/bios
EnvironmentFile=-/etc/default/fty__fty-metric-composite__%i.conf
ExecStart=@prefix@/bin/fty-metric-composite /var/lib/bios/composite-metrics/%i.cfg
Restart=always

[Install]
WantedBy=bios.target
